# deploy/tradingbot.service
[Unit]
Description=TradingBot (APScheduler)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/main.py
Restart=always
RestartSec=10
EnvironmentFile=-__APP_DIR__/.env

# Hardening (optional; remove if they cause issues with ccxt/network)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
